# docs/systemd/betreuungskalender-backup.timer
[Unit]
Description=Daily Betreuungskalender SQLite backup

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=30m

[Install]
WantedBy=timers.target
